# test/fixtures/systemd/app-cleanup-job.service
[Unit]
Description=Run App CleanupJob

[Service]
Type=oneshot
ExecStart=/usr/local/bin/kamal-rails-runner app "CleanupJob.perform_now"
